# .config/systemd/user/kanshi.service
[Unit]
Description=Dynamic display configuration for Wayland
Documentation=man:kanshi(1)
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kanshi
ExecReload=/usr/bin/kanshictl reload

[Install]
WantedBy=sway-session.target
